# netboot.service
[Unit]
Description=Startup Script

[Service]
ExecStart=/usr/local/bin/netboot.sh

[Install]
WantedBy=multi-user.target
